# cadpilot.service
[Unit]
Description=CadPilot prototype
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/cadpilot
ExecStart=/usr/bin/python3 -m http.server 8088 --bind 127.0.0.1
Restart=always
RestartSec=3

[Install]
WantedBy=default.target
